# deploy/container/worker-ami/dirextalk-worker-installer.service
[Unit]
Description=Dirextalk approval-bound Worker root installer
Documentation=https://github.com/YingSuiAI/dirextalk-agent
Requires=dirextalk-worker-installer.socket
After=dirextalk-worker-installer.socket
ConditionFileIsExecutable=/usr/local/bin/dirextalk-worker-installer
ConditionPathExists=/etc/dirextalk-installer/trust.cbor

[Service]
Type=exec
User=root
Group=root
ExecStart=/usr/local/bin/dirextalk-worker-installer serve
UMask=0077
StateDirectory=dirextalk-installer
StateDirectoryMode=0700

# This is the approved privileged boundary of an exclusive Worker VM. The
# daemon still accepts only digest-locked, device-approved commands, but child
# processes may install packages, write system service files, mount approved
# volumes, and use the approved network. Do not apply filesystem, network,
# namespace, device, capability, or syscall sandboxes here.
NoNewPrivileges=yes
LockPersonality=yes
PrivateTmp=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
RemoveIPC=yes
RestrictRealtime=yes
SystemCallArchitectures=native
