# deploy/container/worker-ami/dirextalk-worker-installer.service
[Unit]
Description=Dirextalk approval-bound Worker installer verifier
Documentation=https://github.com/YingSuiAI/dirextalk-agent
Requires=dirextalk-worker-installer.socket
After=dirextalk-worker-installer.socket
ConditionFileIsExecutable=/usr/local/bin/dirextalk-worker-installer
ConditionPathExists=/etc/dirextalk-installer/approval-public-key
ConditionPathExists=/etc/dirextalk-installer/binding.cbor

[Service]
Type=exec
User=root
Group=root
ExecStart=/usr/local/bin/dirextalk-worker-installer
Environment=DIREXTALK_INSTALLER_APPROVAL_PUBLIC_KEY_FILE=/etc/dirextalk-installer/approval-public-key
Environment=DIREXTALK_INSTALLER_BINDING_FILE=/etc/dirextalk-installer/binding.cbor
UMask=0077

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SocketBindDeny=any
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
IPAddressDeny=any
