# deploy/container/worker-ami/dirextalk-worker-installer-bootstrap.service
[Unit]
Description=Dirextalk per-deployment Worker installer trust bootstrap
Documentation=https://github.com/YingSuiAI/dirextalk-agent
Wants=network-online.target
After=network-online.target
Before=dirextalk-cloud-worker.service dirextalk-worker-installer.socket
ConditionPathIsExecutable=/usr/local/bin/dirextalk-worker-installer

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/dirextalk-worker-installer bootstrap
RemainAfterExit=yes
Restart=on-failure
RestartSec=5s
UMask=0077

NoNewPrivileges=yes
# The bootstrap mounts approval-bound EBS filesystems into the host mount
# namespace before any Worker or installer socket can start. Filesystem
# namespace hardening (ProtectSystem/PrivateDevices/PrivateTmp) would hide
# those mounts from the services that consume them, so this one-shot keeps the
# host namespace and retains only the capabilities required for block mounts
# and root-owned path/fstab updates.
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_OVERRIDE
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
