# deploy/container/worker-ami/dirextalk-cloud-worker.service
[Unit]
Description=Dirextalk exclusive-VM Cloud Worker
Documentation=https://github.com/YingSuiAI/dirextalk-agent
Wants=network-online.target
Requires=dirextalk-worker-installer-bootstrap.service
After=network-online.target dirextalk-worker-installer-bootstrap.service
ConditionPathIsExecutable=/usr/local/bin/dirextalk-cloud-worker
ConditionPathIsReadWrite=/var/lib/dirextalk-worker

[Service]
Type=exec
User=dirextalk-worker
Group=dirextalk-worker
WorkingDirectory=/var/lib/dirextalk-worker
StateDirectory=dirextalk-worker
StateDirectoryMode=0700
Environment=DIREXTALK_WORKER_BINARY_SHA256_FILE=/usr/local/share/dirextalk-worker/dirextalk-cloud-worker.sha256
Environment=DIREXTALK_WORKER_TLS_CA_FILE=/etc/ssl/certs/ca-certificates.crt
ExecStart=/usr/local/bin/dirextalk-cloud-worker
Restart=on-failure
RestartSec=10s
UMask=0077

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SocketBindDeny=any
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
